# debian/prerm
#!/bin/bash

sudo systemctl disable docker.onlyoffice.service
